# daemon-trigger-skygear/doorlock-skygear.service
[Unit]
Description=Doorlock Skygear Trigger Daemon

[Service]
ExecStart=/usr/bin/node /home/oursky/skygear

[Install]
WantedBy=multi-user.target
