# daemon-trigger-skygear/doorlock-skygear@.service
[Unit]
Description=Doorlock Skygear Trigger Daemon

[Service]
ExecStart=/usr/bin/node /home/oursky/doorlock/daemon-trigger-skygear %i
Restart=on-failure

[Install]
WantedBy=multi-user.target
