# daemon-trigger-ble/doorlock-ble@.service
[Unit]
Description=Doorlock BLE Advertising Daemon

[Service]
ExecStart=/usr/bin/node /home/oursky/ble %i

[Install]
WantedBy=multi-user.target
